# image/files/automount/usb-automount@.service
[Unit]
Description=USB automount for %I
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb-automount.sh %I
